# reporting-gunicorn.conf
[program:reporting-gunicorn]
command = gunicorn --bind 0.0.0.0:8080 reportingauto.wsgi:application
user = app
directory = /home/app/web/reportingauto
autostart = true
autorestart = true